# led_monitor.service
[Unit]
Description=LED Status Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/led_monitor_script.sh
Restart=on-failure
RestartSec=5
User=root

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

